# ansible/roles/patroni/files/patroni.service
[Unit]
Description=Patroni PostgreSQL HA Cluster Node
After=network.target

[Service]
Type=simple
User=postgres
ExecStart=/usr/local/bin/patroni /etc/patroni/postgres.yml
Restart=on-failure
RestartSec=5s
LimitNOFILE=262144

[Install]
WantedBy=multi-user.target
